# .history/startvm_20210908112410.service
[unit]
Description=Start OpenMRS vm
After=network.target vboxdrv.service
Before=runlevel2.target

[Service]
User=hms
Group=vboxusers
Type=forking
Restart=no 
TimeoutSec=5min
KillMode=process
RemainAfterExit=yes